# dc_event_detection.service
[Unit]
Description=DC Event Detection Service
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=simple
ProtectHome=true

ExecStart=/usr/local/bin/event_detect /etc/event_detect.conf

Restart=on-failure
RestartSec=5s
# Optional: Consider StartLimitIntervalSec=, StartLimitBurst= for rate limiting.

RuntimeDirectory=event_detect
RuntimeDirectoryMode=0755 # Good: Allows other users/groups to access the directory.

# Recommendation: Add User= and Group= for security. HIGHLY RECOMMENDED.
# Create a dedicated system user/group (e.g., 'event-detect') first:
#   sudo groupadd --system event-detect
#   sudo useradd --system -g event-detect -d / -s /sbin/nologin event-detect
# Add the user to groups needed for device access (check exact names on openSUSE):
#   sudo usermod -aG input event-detect
#   sudo usermod -aG tty event-detect # Or maybe 'dialout'? Check TTY permissions.
#User=event-detect
#Group=event-detect

# Recommendation: Consider additional security hardening if possible.
# PrivateTmp=yes # Usually safe and recommended.
# NoNewPrivileges=yes # Good security practice.
# ProtectSystem=full # Maybe 'strict' if it doesn't break device access needed.
                     # Start with 'full' or leave it off initially.
# PrivateDevices=yes # NO! This service needs device access.

[Install]
WantedBy=multi-user.target # Correct for a system service started on boot.
