# dc_event_detection.service.in
[Unit]
Description=DC Event Detection Service
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=simple
ProtectHome=true

ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/event_detect /etc/event_detect.conf

Restart=on-failure
RestartSec=5s

RuntimeDirectory=event_detect
RuntimeDirectoryMode=0755

User=event_detect
Group=event_detect

# Filesystem Access Control
ProtectHome=true
ProtectSystem=strict
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Process Execution Control
NoNewPrivileges=yes

# Network Access Control
PrivateNetwork=yes

[Install]
WantedBy=multi-user.target
