# dc_idle_detection.service.in
[Unit]
Description=BOINC Idle Detection Service
After=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/idle_detect_wrapper.sh
Restart=on-failure
RestartSec=5s
ProtectSystem=strict
ReadWritePaths=/run/event_detect
NoNewPrivileges=yes

[Install]
WantedBy=graphical-session.target
