# debian.prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ] || [ "$1" = "upgrade" ] || [ "$1" = "deconfigure" ]; then
    # Use dh-systemd helpers (called via debian/rules)
    # Or manually disable/stop here:
    # echo "Disabling/Stopping systemd service dc_event_detection.service..."
    # systemctl disable --now dc_event_detection.service || true
    : # No-op if relying on dh helpers
fi
# Include debhelper snippets (important!)
#DEBHELPER#
exit 0
